# pi/launchd/com.peterg17.obsidian-index-cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.peterg17.obsidian-index-cleanup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/peterg17/Documents/notes/.pi/scripts/process-index.sh</string>
  </array>

  <!-- Runs daily at 8:30 PM local time. Edit Hour/Minute if desired. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>20</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/peterg17/Documents/notes</string>

  <key>StandardOutPath</key>
  <string>/Users/peterg17/Documents/notes/.pi/logs/index-cleanup.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/peterg17/Documents/notes/.pi/logs/index-cleanup.err.log</string>
</dict>
</plist>
