# pi/launchd/com.peterg17.obsidian-ticket-kanban-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.peterg17.obsidian-ticket-kanban-refresh</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/peterg17/Documents/notes/.pi/scripts/rebuild-obsidian-kanban.sh</string>
  </array>

  <!-- Re-render the generated Obsidian Kanban board every 10 minutes. -->
  <key>StartInterval</key>
  <integer>600</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>WorkingDirectory</key>
  <string>/Users/peterg17/Documents/notes</string>

  <key>StandardOutPath</key>
  <string>/Users/peterg17/Documents/notes/.pi/logs/kanban-refresh.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/peterg17/Documents/notes/.pi/logs/kanban-refresh.err.log</string>
</dict>
</plist>
